# private_dot_config/systemd/user/gammastep.service.d/override.conf
[Unit]
PartOf=
PartOf=sway-session.target
After=
After=sway-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Install]
WantedBy=
WantedBy=sway-session.target
